# snowflake/templates/default/snowflaked.conf.erb
description "snowflake id generator service"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

limit nofile 8192 8192

script
    exec sudo -u snowflake /usr/bin/java -cp /usr/share/snowflake/libs/*:/usr/share/snowflake/* com.twitter.service.snowflake.SnowflakeServer -f /etc/snowflake/snowflake.scala
end script
